# deploy/native/go2rtc.service
[Unit]
Description=go2rtc streaming engine (camera-hub-personal)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__PREFIX__/go2rtc
ExecStart=__PREFIX__/go2rtc/go2rtc -config __PREFIX__/go2rtc/go2rtc.yaml
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
